# seps_web/rakkun-api.service
# ─────────────────────────────────────────────────────────────
# Systemd service — Rakkun FastAPI (Uvicorn)
# Copiar a: /etc/systemd/system/rakkun-api.service
# Luego:
#   systemctl daemon-reload
#   systemctl enable rakkun-api
#   systemctl start rakkun-api
# ─────────────────────────────────────────────────────────────

[Unit]
Description=Rakkun API — FastAPI / Uvicorn
After=network.target postgresql.service

[Service]
User=root
WorkingDirectory=/var/www/intellicoop/seps_web
EnvironmentFile=/var/www/intellicoop/seps_web/.env

ExecStart=/var/www/intellicoop/venv/bin/uvicorn \
    rakkun_api:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 1

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
